# deploy/systemd/fhirsmith-prod.service
[Unit]
Description=FHIRsmith (FHIRTX) — production
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/fhirsmith-prod
Environment=NODE_ENV=production
Environment=PORT=3002
Environment=FHIRSMITH_DATA_DIR=/var/lib/fhirsmith-prod
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=10s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
