# deploy/systemd/fhirsmith-staging.service
[Unit]
Description=FHIRsmith (FHIRTX) — staging
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/fhirsmith-staging
Environment=NODE_ENV=staging
Environment=PORT=3001
Environment=FHIRSMITH_DATA_DIR=/var/lib/fhirsmith-staging
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=10s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
